# deploy/systemd/catty-app.service
[Unit]
Description=Catty reminders app
After=network.target

[Service]
Type=simple
User=catty
Group=catty
WorkingDirectory=/opt/catty/app
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/catty/app/.env
ExecStart=/opt/catty/app/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8181
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target